# roles/brewmaster-frontend/templates/nodejs_unit_file.service.j2
[Unit]
Description="Node.js frontend"
After=network.target

[Service]
WorkingDirectory={{ repo_absolute_path.stdout }}
ExecStart=npm run-script "{{ npm['run_script'] }}"
