# roles/brewmaster-frontend/templates/nodejs_unit_file_file.service.j2
[Unit]
Description="Node.js frontend"
After=network.target

[Service]
ExecStart=npm run-script "{{ npm['run_script'] }}"
