# roles/brewmaster-backend/templates/flask_unit_file.service.j2
[Unit]
Description="Backend Flask server application"
After=network.target

[Service]
Environment=FLASK_APP={{ repo_absolute_path['stdout'] }}/application
ExecStart=flask run
